# deploy/daily-summary-web.service
# systemd USER service for the dashboard (FastAPI/Uvicorn). Installed as a *user* service (no sudo);
# linger keeps it up with no active login. Install:
#   loginctl enable-linger "$USER"        # once, if not already enabled
#   mkdir -p ~/.config/systemd/user
#   cp deploy/daily-summary-web.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now daily-summary-web
#   sudo ufw allow from 192.168.1.0/24 to any port 8080 proto tcp   # open the LAN port (one-time)
# Then Avigail opens http://<mini-pc-ip>:8080  (no auth — LAN only).
[Unit]
Description=daily-summary dashboard (FastAPI/Uvicorn, LAN-only)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=%h/daily-summary
ExecStart=%h/.local/bin/uv run uvicorn digest_web.app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
